# systemd/vnc-tunnel@.service
[Unit]
Description=VNC Tunnel Service for %i
Documentation=https://github.com/SGT-Cho/cloudflared-vnc-tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/vnc-tunnel %i
Restart=always
RestartSec=10
User=vnc-tunnel
Group=vnc-tunnel

# Environment
Environment="HOME=/var/lib/vnc-tunnel"
Environment="USER=vnc-tunnel"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/vnc-tunnel
# Allow cloudflared to store auth tokens
ReadWritePaths=/var/lib/vnc-tunnel
StateDirectory=vnc-tunnel

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vnc-tunnel-%i

[Install]
WantedBy=multi-user.target